# cron.d/cron_conf
*/15 * * * * root /etc/scripts/incron_config.sh
0 0 * * * root /etc/scripts/remote_backup.sh
0 0 * * * root /etc/scripts/renew_cert.sh
